# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:oidentd]
redirect_stderr=true
command=oidentd -i -u nobody -g nogroup
